# ansible/roles/github_runner/templates/github-runner.service.j2
[Unit]
Description=GitHub Actions Runner ({{ github_runner_name }})
After=network.target

[Service]
Type=simple
User={{ github_runner_user }}
Group={{ github_runner_user }}
WorkingDirectory={{ github_runner_dir }}
ExecStart={{ github_runner_dir }}/run.sh
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=5min

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths={{ github_runner_dir }}
PrivateTmp=true

# Environment
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1
Environment=RUNNER_ALLOW_RUNASROOT=0

[Install]
WantedBy=multi-user.target
